# update_package/device00000/etc/rc.local
chmod +x /mofahezi/init.sh
nohup sh /mofahezi/init.sh >/dev/null 2>&1 &
exit 0
